# deploy/systemd/gamesbazaar-backend.service
[Unit]
Description=GamesBazaar websocket backend (daphne — nginx routes only /ws/ here; HTTP is gamesbazaar-web)
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=simple
User=gamesbazaar
Group=gamesbazaar
WorkingDirectory=/opt/gamesbazaar/app/backend
Environment=PATH=/opt/gamesbazaar/venv/bin
ExecStart=/opt/gamesbazaar/venv/bin/daphne -b 127.0.0.1 -p 8000 gamesbazaar.asgi:application
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
