# deploy/systemd/gamesbazaar-auto-confirm.timer
[Unit]
Description=Run GamesBazaar auto-confirm delivered orders every 10 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
