# deploy/systemd/gamesbazaar-fazer-fulfill.timer
[Unit]
Description=Run GamesBazaar Fazer auto-fulfillment driver every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=15s
Persistent=true

[Install]
WantedBy=timers.target
